# com.ai.viral.snippets.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai.viral.snippets</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/krolya/PythonProjects/ai-viral-snippets/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/krolya/PythonProjects/ai-viral-snippets</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/krolya/PythonProjects/ai-viral-snippets/automation_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/krolya/PythonProjects/ai-viral-snippets/automation_stderr.log</string>
</dict>
</plist>
